# deploy/greenhouse-contractor-caller-call.service
[Unit]
Description=Greenhouse Contractor Caller outbound run

[Service]
Type=oneshot
WorkingDirectory=/opt/greenhouse-contractor-caller
EnvironmentFile=/etc/greenhouse-contractor-caller/env
ExecStart=/opt/greenhouse-contractor-caller/.venv/bin/python -m app.caller

